# deploy/systemd/cloudflared-replica.service
# Template: paths below are for the active gate. gate_provision.sh
# rewrites User/Group and paths when installing on replicas.
[Unit]
Description=cloudflared tunnel replica (nucleus-lab)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=irongate
Group=irongate
EnvironmentFile=-/etc/projectnucleus/gate.env
ExecStart=${GATE_HOME}/bin/cloudflared tunnel --config ${GATE_HOME}/.cloudflared/config.yml run nucleus-lab
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
