# deploy/systemd/sporeprint-rebuild.service
# Template: paths below are for the active gate. gate_provision.sh
# rewrites User/Group and paths when installing on replicas.
[Unit]
Description=sporePrint rebuild (git pull + zola build)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=irongate
Group=irongate
EnvironmentFile=-/etc/projectnucleus/gate.env
WorkingDirectory=${GATE_HOME}/Development/ecoPrimals/sporeGarden/projectNUCLEUS
ExecStart=/usr/bin/bash deploy/sporeprint_local.sh build
TimeoutStartSec=120
