# deploy/systemd/sporeprint-local.service
# Template: paths below are for the active gate. gate_provision.sh
# rewrites User/Group and paths when installing on replicas.
[Unit]
Description=sporePrint local preview (dev/preview only — primals.eco served by GitHub Pages CDN)
After=network.target
Documentation=https://github.com/sporeGarden/projectNUCLEUS

[Service]
Type=simple
User=irongate
Group=irongate
EnvironmentFile=-/etc/projectnucleus/gate.env
WorkingDirectory=${GATE_HOME}/Development/ecoPrimals/sporeGarden/projectNUCLEUS
ExecStart=/usr/bin/bash deploy/sporeprint_local.sh serve
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
